# ansible/roles/dataiku-dss/files/systemd/configure-dataiku.service
[Unit]
Description=A service that configures dataiku
After=dataiku.service
ConditionPathExists=!/var/run/configure-dataiku.lock

[Service]
Type=oneshot
Environment=PYTHONPATH=/opt/dataiku-python
ExecStartPre=/bin/touch /var/run/configure-dataiku.lock
WorkingDirectory=/opt/dataiku-python/dataiku-configurator
ExecStart=/bin/python3 /opt/dataiku-python/dataiku-configurator/main.py
StandardOutput=journal
SyslogIdentifier=configure-dataiku

[Install]
WantedBy=multi-user.target
